# cloudbot/incubator/openclaw-extensions/devops-sre/systemd/devops-healthcheck.service
[Unit]
Description=DevOps daily health check
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/devops-sre
ExecStart=/opt/devops-sre/checks/daily_check.sh --services "nginx docker postgresql redis-server"
