# cloudbot/incubator/openclaw-extensions/devops-sre/systemd/devops-weekly-audit.timer
[Unit]
Description=Run DevOps weekly audit on Monday

[Timer]
OnCalendar=Mon *-*-* 09:30:00 Europe/Moscow
Persistent=true
Unit=devops-weekly-audit.service

[Install]
WantedBy=timers.target
