# cloudbot/incubator/openclaw-extensions/devops-sre/systemd/devops-weekly-audit.service
[Unit]
Description=DevOps weekly reliability audit
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/devops-sre
ExecStart=/opt/devops-sre/checks/weekly_check.sh --services "nginx docker postgresql redis-server"
