# tooling/claude-state-sync/launchd/com.user.state-drift-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.state-drift-watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/pro2kuror/.claude/hooks/state-drift-check.sh</string>
    <string>full</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/pro2kuror/.claude/drift-watchdog.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pro2kuror/.claude/drift-watchdog.log</string>
</dict>
</plist>
