# deploy/cron/empty_companies_weekly.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
TZ=Europe/Moscow

# «Пустые компании» — weekly refresh в низкое окно 02:00-06:00 МСК.
# Лог wrapper'а: ~/work-crm-enrich/belberry/bitrix24/logs/empty_companies_cron.log
30 3 * * 1 root WORKTREE=/root/work-crm-enrich MODULE_DIR=/opt/cloudbot-runtime/larisa/empty-companies-refresh/belberry/bitrix24/crm_deal_merge VENV=/opt/cloudbot-runtime/larisa/sales-kpi-dashboard/.venv BITRIX_STATE_PATH=/opt/openclaw/state/bitrix_app/install.latest.json BITRIX_SYNC_SCRIPT=/opt/openclaw/repos/vibecoding/shared/scripts/bitrix-sync-state.sh GOOGLE_SA_KEY=/opt/openclaw/secrets/finance-director-sheets-903611b799c3.json /opt/cloudbot-runtime/larisa/empty-companies-refresh/belberry/bitrix24/crm_deal_merge/scripts/cron_refresh.sh >> /root/work-crm-enrich/belberry/bitrix24/logs/empty_companies_cron.log 2>&1
