# infra/templates/sing-box.service
[Unit]
Description=Happ VPN Data Plane (sing-box)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=happvpn
Group=happvpn
Environment=TZ=Europe/Moscow
ExecStart=/usr/local/bin/sing-box run -c /etc/sing-box/config.json
Restart=always
RestartSec=3
LimitNOFILE=1048576
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
